# deploy/ansible-web.ini
[program:ansible-web]
command=/srv/AnsibleConfigurationWeb/venv/bin/gunicorn_start
user=ec2-user
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/AnsibleConfigurationWeb/logs/gunicorn-error.log
